# .config/systemd/user/mbsync.timer
[Unit]
Description=Run mbsync every 15 minutes

[Timer]
OnBootSec=5m
OnUnitActiveSec=15m
Unit=mbsync.service

[Install]
WantedBy=timers.target
