# .config/systemd/user/mbsync.service
[Unit]
Description=Mailbox synchronization

[Service]
Type=oneshot
ExecStart=%h/.local/bin/mail-sync
